# deploy/ourtab.cron
20 3 * * * root cd /home/ubuntu/apps/ourtab && /bin/bash deploy/backup.sh >> /home/ubuntu/apps/ourtab/backups/backup.log 2>&1
